# systemd/hls-gallery-media-permissions.path
[Unit]
Description=Watch HLS Video Gallery (@@INSTANCE_ID@@) for new media uploads

[Path]
PathChanged=@@DOCUMENT_ROOT@@/media
Unit=hls-gallery-@@INSTANCE_ID@@-media-permissions.service

[Install]
WantedBy=multi-user.target
