# systemd/hls-gallery-media-permissions.timer
[Unit]
Description=Periodically repair HLS Video Gallery (@@INSTANCE_ID@@) media permissions

[Timer]
OnBootSec=20s
OnUnitInactiveSec=1min
AccuracySec=5s
Persistent=true
Unit=hls-gallery-@@INSTANCE_ID@@-media-permissions.service

[Install]
WantedBy=timers.target
